# server/src/dist/rhel7-linux-x86_64/systemd/elasticsearch@.service
[Unit]
Description=Elasticsearch server %I
After=network.target

[Service]
Type=simple
User=es
EnvironmentFile=-/etc/sysconfig/elasticsearch
ExecStart=/usr/bin/elasticsearch %i
Restart=on-abort

[Install]
WantedBy=multi-user.target
