# server/src/dist/rhel7-x86_64/rpm/java-11/default/etc/systemd/system/elasticsearch.service
[Unit]
Description=Elasticsearch server
After=network.target

[Service]
Type=simple
User=es
EnvironmentFile=-/etc/sysconfig/elasticsearch-default
ExecStart=/opt/elasticsearch/bin/elasticsearch-default
Restart=on-failure
LimitNOFILE=65536
LimitNPROC=4096
LimitMEMLOCK=unlimited

[Install]
WantedBy=multi-user.target
